# launchagents/com.leosaquetto.codexusage.antigravity.autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.leosaquetto.codexusage.antigravity.autoupdate</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/leosaquetto/Developer/GitHub/codex-usage/scripts/update-antigravity-usage-auto.mjs</string>
		<string>--commit</string>
		<string>--push</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/leosaquetto/Developer/GitHub/codex-usage/.local/logs/antigravity-usage-launchd.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/leosaquetto/Developer/GitHub/codex-usage/.local/logs/antigravity-usage-launchd.out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Minute</key>
		<integer>23</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/leosaquetto/Developer/GitHub/codex-usage</string>
</dict>
</plist>
